# deploy/cryptobot-live.service
[Unit]
Description=CryptoBot Live Service
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/CryptoBot
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/%i/CryptoBot/.venv/bin:/usr/local/bin:/usr/bin
EnvironmentFile=/home/%i/CryptoBot/.env
ExecStart=/home/%i/CryptoBot/.venv/bin/python -m cryptobot.cli.live --config /home/%i/CryptoBot/configs/live.frugal.yaml --provider ccxt --exchange binance
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
